# usr/lib/systemd/system/vpn-ping@.service
[Unit]
Description=VPN ping test for location : %I

[Service]
Type=simple
User=root
ExecStart=/etc/openvpn/pia/vpn-ping.sh 8.8.8.8 %I
